# scripts/crontab
# Cron jobs do Album Copa Samsung 2026
# Edite este arquivo, faça git push → git pull no servidor aplica automaticamente.
# Formato: minuto hora dia mês dia-semana usuário comando

# Distribuição diária de pacotes — 12:00 BRT (15:00 UTC), seg–sex
0 15 * * 1-5 root /var/www/album-figurinhas/scripts/distribuir.sh

# Pacote bônus de fim de semana — 14:00 BRT (17:00 UTC), sáb–dom, 5 figurinhas
0 17 * * 0,6 root /var/www/album-figurinhas/scripts/fim-de-semana.sh
